# sites/_template_wordpress/etc/crontab
# m h  dom mon dow  user  command
* * * * * application wget -q -O - http://localhost/wp-cron.php?doing_wp_cron >/dev/null 2>&1
